# bikebinder.conf
[program:bikebinder]
directory=/BikeBinder
command=rails server
autorestart=true
